# etc/systemd/system/securonis-timezone-spoofer.service
[Unit]
Description=Securonis - Timezone Spoofer
After=network.target
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/securonis-timezone-spoofer
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
