# etc/systemd/system/securonis-cleanup-temp.service
[Unit]
Description=Securonis - Clean Temporary Files on Shutdown
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target
Conflicts=reboot.target

[Service]
Type=oneshot
ExecStart=/usr/bin/securonis-cleanup-temp
TimeoutStartSec=30
RemainAfterExit=no

[Install]
WantedBy=shutdown.target
